# files/etc/crontabs/root
15 01 * * 4 bash /root/scripts/update.sh 0
