# files/etc/rc.local
#!/bin/sh -e                                                                                                                                                   
#
# rc.local
#
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

cd /root/scripts/

./start.sh &

exit 0
